# fs/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin hackeitos --noclear %I - linux
# ExecStart=-/sbin/agetty --noclear /usr/bin/ly - linux # No autologin, cool login screen